# dist/supervisor.conf
[program:ceitba]
directory = /home/sebikul/CEITBA
command = /home/sebikul/CEITBA/dist/runserver.sh
user = sebikul
stdout_logfile = /home/sebikul/CEITBA/logs/supervisor.log
redirect_stderr = true
environment = LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8,LC_LANG=en_US.UTF-8